# ops/launchd/com.openai.codex.repe-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.openai.codex.repe-bridge</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/Cellar/python@3.11/3.11.12/Frameworks/Python.framework/Versions/3.11/bin/python3.11</string>
		<string>/Users/paulmalmquist/VSCodeProjects/BusinessMachine/Consulting_app/scripts/repe_codex_os_waker.py</string>
		<string>--bridge</string>
		<string>/Users/paulmalmquist/VSCodeProjects/BusinessMachine/Consulting_app/verification/loop_state/repe_supervised_bridge.json</string>
		<string>--state</string>
		<string>/Users/paulmalmquist/VSCodeProjects/BusinessMachine/Consulting_app/verification/loop_state/repe_codex_os_waker_state.json</string>
		<string>--automation-id</string>
		<string>claude-repe-watch</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/paulmalmquist/VSCodeProjects/BusinessMachine/Consulting_app/verification/loop_state/logs/repe_codex_os_waker.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/paulmalmquist/VSCodeProjects/BusinessMachine/Consulting_app/verification/loop_state/logs/repe_codex_os_waker.out.log</string>
	<key>StartInterval</key>
	<integer>60</integer>
	<key>WatchPaths</key>
	<array>
		<string>/Users/paulmalmquist/VSCodeProjects/BusinessMachine/Consulting_app/verification/loop_state/repe_supervised_bridge.json</string>
	</array>
</dict>
</plist>
